# 99-screenconfig.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/etc/udev/scripts/screenconfig.sh"
